# sendexternalip.conf
start on startup
task
exec /usr/lib/sendexternalip.bash
